# systemd/ripple.service
[Unit]
Description=Firebolt reference implementation
Requires=wpeframework.service
After=wpeframework.service

[Service]
Type=notify
Environment="RUST_LOG=debug"
EnvironmentFile=-/opt/ripple.conf
ExecStart=/usr/bin/ripple
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
